# misc/crontab
0 */2 * * * /home/ori/tepco/bin/mirror.sh > /home/ori/tepco/log 2>&1
